# docker/s6-overlay/s6-rc.d/proxyos/finish
#!/command/with-contenv sh
# Called by s6 after the Node process exits. Nothing extra needed.
exit 0
