# docker/s6-overlay/s6-rc.d/caddy/run
#!/bin/sh
exec /usr/local/bin/caddy run --config /config/caddy/base-config.json
